# janitor.service
[Unit]
Description=Janitor

[Service]
ExecStart=~/.cargo/bin/janitor-service
Restart=always

[Install]
WantedBy=default.target