# janitor.service
[Unit]
Description=Janitor

[Service]
ExecStart=/home/sharun/.cargo/bin/janitor-service
Restart=always

[Install]
WantedBy=graphical.target